# sensor-web.service
[Unit]
Description=Pi Temperature Sensor Web Interface
After=network.target

[Service]
User=pi
WorkingDirectory=/home/pi/pi_temp_sensor
ExecStart=/home/pi/pi_temp_sensor/venv/bin/python3 /home/pi/pi_temp_sensor/app.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
